# supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
pidfile=/var/run/supervisord.pid
loglevel=info

[program:rag_server]
command=python3 rag_server.py --port 5080
autostart=true
autorestart=true
priority=100
startsecs=5
stderr_logfile=/app/logs/rag.err.log
stdout_logfile=/app/logs/rag.out.log

[program:moduleflow]
command=python3 moduleflow/serve.py
autostart=true
autorestart=true
priority=200
startsecs=3
stderr_logfile=/app/logs/moduleflow.err.log
stdout_logfile=/app/logs/moduleflow.out.log

[program:main_app]
command=python3 -m uvicorn main:app --host 0.0.0.0 --port 5069 --log-level info
autostart=true
autorestart=true
priority=300
startsecs=5
stderr_logfile=/app/logs/main.err.log
stdout_logfile=/app/logs/main.out.log
